# butler-log-import.service
[Unit]
Description=Butler — импорт попыток подключений из journald
After=network.target

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
